# offlineimap/offlineimap-oneshot.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=10m

[Install]
WantedBy=default.target
